# s6/php-fpm7/run
#!/bin/sh

exec /usr/sbin/php-fpm7 -R -F --fpm-config /etc/php7/php-fpm.conf
